# Lesson 6 System initialization. Systemd and SysV./2/spawn-fcgi.socket
[Unit]
Description=spawn-fcgi Listen Socket

[Socket]
ListenStream=/var/run/spawn-fcgi.sock

[Install]
WantedBy=sockets.target

