# Lesson 6 System initialization. Systemd and SysV./1/myservice.service
Unit]
Description=Logfile word finder

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/myserviceconfig
ExecStart=/bin/bash /root/myservice.sh

